# scripts/remap-keys.service
[Unit]
Description=Remap f4 and f5 using evsieve
[Service]
ExecStart=/usr/local/sbin/remap.sh
[Install]
WantedBy=multi-user.target
